# crontab
0 1 15 1,3,6,9,12 * cd /home/andy/sh && sh BackupImmich.sh &
0 1 17 1,3,6,9,12 * cd /home/andy/sh && sh RestoreBackUp.sh &